# deploy/systemd/knowledgehub-zotero-sync.timer
[Unit]
Description=Run KnowledgeHub Zotero sync hourly

[Timer]
OnCalendar=hourly
Persistent=true
AccuracySec=1min
RandomizedDelaySec=10min
Unit=knowledgehub-zotero-sync.service

[Install]
WantedBy=timers.target
