# deploy/systemd/knowledgehub-zotero-sync.service
[Unit]
Description=KnowledgeHub Zotero incremental sync
Documentation=file:///home/lengmo/KnowledgeHub/docs/sources/zotero.md
Wants=network-online.target
Requires=knowledgehub-zotero-cache-refresh.service
After=network-online.target knowledgehub-zotero-cache-refresh.service

[Service]
Type=oneshot
User=lengmo
WorkingDirectory=/home/lengmo/KnowledgeHub
EnvironmentFile=/etc/knowledgehub/zotero.env
ExecStart=/home/lengmo/anaconda3/bin/conda run -n rag --no-capture-output knowledgehub --config /etc/knowledgehub/zotero.yaml zotero sync --once
TimeoutStartSec=30min
UMask=0077

# The source resolver sees an immutable attachment mirror. Its prerequisite
# paginated WebDAV refresh is the only writer and exits before this unit starts.
ProtectSystem=strict
ProtectHome=read-only
ReadOnlyPaths=/home/lengmo/KnowledgeHub /home/lengmo/anaconda3
ReadOnlyPaths=/data/KnowledgeHub/zotero_cache
ReadWritePaths=/data/KnowledgeHub/zotero
PrivateTmp=true
PrivateDevices=true
NoNewPrivileges=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
LockPersonality=true
RestrictRealtime=true
RestrictSUIDSGID=true
SystemCallArchitectures=native

# stdout/stderr, including the structured run summary, remain in the journal.
StandardOutput=journal
StandardError=journal
SyslogIdentifier=knowledgehub-zotero-sync
