# deploy/systemd/knowledgehub-zotero-rag-incremental.service
[Unit]
Description=KnowledgeHub Zotero source and RAG incremental ingest
Requires=knowledgehub-zotero-sync.service knowledgehub-rag-core.service
After=knowledgehub-zotero-sync.service knowledgehub-rag-core.service network-online.target

[Service]
Type=oneshot
User=lengmo
WorkingDirectory=/home/lengmo/KnowledgeHub
EnvironmentFile=/etc/knowledgehub/rag.env
ExecStart=/usr/local/libexec/knowledgehub-rag-incremental-with-retries
TimeoutStartSec=30h
TimeoutStopSec=5min
UMask=0077
ProtectSystem=strict
ProtectHome=read-only
ReadOnlyPaths=/home/lengmo/KnowledgeHub /home/lengmo/anaconda3 /data/KnowledgeHub/zotero
ReadWritePaths=/data/KnowledgeHub/rag/zotero /data/KnowledgeHub/model-cache
NoNewPrivileges=true
PrivateTmp=true
StandardOutput=journal
StandardError=journal
SyslogIdentifier=knowledgehub-zotero-rag
