# deploy/systemd/knowledgehub-zotero-cache-refresh.service
[Unit]
Description=Refresh the local KnowledgeHub Zotero attachment mirror
Documentation=file:///home/lengmo/KnowledgeHub/docs/sources/zotero.md
Wants=network-online.target
After=network-online.target
Before=knowledgehub-zotero-sync.service

[Service]
Type=oneshot
User=lengmo
ExecStartPre=/usr/bin/mkdir -p /data/KnowledgeHub/zotero_cache
ExecStart=/usr/bin/rclone sync nutstore:/zotero /data/KnowledgeHub/zotero_cache --transfers 1 --checkers 2 --tpslimit 1 --tpslimit-burst 1 --low-level-retries 3 --retries 3 --retries-sleep 60s --delete-after --log-level INFO
TimeoutStartSec=2h
UMask=0077

# Only the disposable local mirror is writable. The remote is the rclone
# source argument, so sync can never delete from nutstore:/zotero.
ProtectSystem=strict
ProtectHome=read-only
ReadOnlyPaths=/home/lengmo/.config/rclone
ReadWritePaths=/data/KnowledgeHub/zotero_cache
PrivateTmp=true
PrivateDevices=true
NoNewPrivileges=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
LockPersonality=true
RestrictRealtime=true
RestrictSUIDSGID=true
SystemCallArchitectures=native

StandardOutput=journal
StandardError=journal
SyslogIdentifier=knowledgehub-zotero-cache-refresh
